# upgrade-rootfs/zynq-rootfs/usr/lib/opkg/info/cgminer.list
/usr/bin/cgminer
/etc/config/cgminer
/usr/bin/cgminer-monitor
/usr/bin/cgminer-api
/etc/init.d/cgminer
/usr/bin/mmupgrade
