# upgrade-rootfs/zynq-rootfs/usr/lib/opkg/info/om-watchdog.list
/etc/init.d/om-watchdog
/sbin/om-watchdog
